# scripts/axon_ros.service
[Unit]
Description=Axon ROS 2 Bringup
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/axon/axon_ws
Environment=ROS_DOMAIN_ID=42
Environment=ROS_LOCALHOST_ONLY=0
ExecStart=/bin/bash -lc "source /opt/ros/jazzy/setup.bash && source /opt/axon/axon_ws/install/setup.bash && ros2 launch axon_bringup robot.launch.py profile:=real"
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
