# scripts/axon_pi.service
[Unit]
Description=Axon Pi Core (ROS 2 Humble)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/pi/Axon-ros
Environment=ROS_DOMAIN_ID=42
ExecStart=/bin/bash -lc 'source /opt/ros/humble/setup.bash && source /home/pi/Axon-ros/install/setup.bash && ros2 launch axon_bringup pi_core.launch.py'
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
